# scripts/mc-idle-shutdown.timer
[Unit]
Description=Minecraft Idle Shutdown Timer
# Use Wants instead of Requires to prevent timer from stopping if minecraft restarts
Wants=minecraft.service
After=minecraft.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=30s
# Keep timer active even if triggered service fails
Persistent=true

[Install]
WantedBy=timers.target
